# webapp-supervisor.conf
[program:webapp]
directory=/test-head
command=sh /usr/local/bin/docker-run.sh
environment=RACK_ENV=production
autostart=true
autorestart=true

